# app-configs/.config/systemd/user/backup-home.service
[Unit]
Description=Borg backup of home to /mnt/backup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-home
Nice=19
IOSchedulingClass=idle
