# app-configs/.config/systemd/user/backup-home.timer
[Unit]
Description=Daily home backup (borg to /mnt/backup)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
